# vnc-image/build/rootfs/etc/supervisor/conf.d/wine-vnc.conf
[group:wine-vnc]
programs=xvfb,wm,pcmanfm,x11vnc

[program:xvfb]
priority=0
command=/usr/local/bin/xvfb.sh
stopsignal=KILL

[program:wm]
priority=5
command=/usr/bin/icewm
environment=DISPLAY=":1",HOME="/root",USER="root"

[program:pcmanfm]
priority=5
user=nobody
directory=/home/nobody
command=/usr/bin/pcmanfm --desktop --profile IceWM
environment=DISPLAY=":1",HOME="/home/nobody",USER="nobody"

[program:x11vnc]
priority=5
command=bash -c "sleep 3 && x11vnc -rfbauth .password2 -forever -xkb -noxrecord -noxdamage -shared -no6 -ncache 0"
environment=DISPLAY=":1",XAUTHLOCALHOSTNAME="localhost"
