# build/rootfs/etc/supervisor/conf.d/tinyproxy.conf
[program:tinyproxy]
priority=10
command=tinyproxy -d -c /etc/tinyproxy/tinyproxy.conf
